# deployment/pitchdock-frontend.service
[Unit]
Description=PitchDock Next.js Frontend Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/emailAutomater/web
ExecStart=/usr/bin/npm run start -- -p 3000
Restart=always
EnvironmentFile=/home/ubuntu/emailAutomater/.env

[Install]
WantedBy=multi-user.target
