# deployment/pitchdock-backend.service
[Unit]
Description=PitchDock FastAPI Backend Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/emailAutomater
ExecStart=/home/ubuntu/emailAutomater/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8000
Restart=always
EnvironmentFile=/home/ubuntu/emailAutomater/.env

[Install]
WantedBy=multi-user.target
